# supervisord.conf
; ==================== supervisord.conf (FIXED) ====================

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/apache2.err.log
stdout_logfile=/var/log/supervisor/apache2.out.log
user=www-data

[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/cron.err.log
stdout_logfile=/var/log/supervisor/cron.out.log

; 🔥 HEALTHCHECK FIX - Simple HTTP check instead of curl
[program:healthcheck]
command=/bin/bash -c "while true; do echo 'Health check running' >> /var/log/health.log; sleep 60; done"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/health.err.log
stdout_logfile=/var/log/supervisor/health.out.log
